# systemd/separation.service
[Unit]
Description=Vocal/Instrumental Separation Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/data/mahdidev/ros/dev
ExecStart=/mnt/data/mahdidev/ros/dev/separation-service/start.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
